# dental_clinic.service
[Unit]
Description=Dental Clinic Flask App
After=network.target

[Service]
User=root
WorkingDirectory=/root/dental_clinic
Environment="PATH=/root/dental_clinic/.venv/bin"
Environment="APP_ENV=production"
Environment="SESSION_COOKIE_SECURE=1"
Environment="WTF_CSRF_SSL_STRICT=0"
ExecStart=/root/dental_clinic/.venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
